# crontab
# SSL証明書自動更新
0 * */2 * * /usr/bin/certbot renew && /usr/bin/systemctl reload nginx