# crontab
# Update ssl certificate at every two months.
0 * */2 * * /usr/bin/certbot renew && /usr/bin/systemctl reload nginx